# deploy/craftbox-update.timer
[Unit]
Description=Daily Craftbox Update Timer (Instance: {{INSTANCE}})
Documentation=https://github.com/raymondjxu/craftbox
Requires=craftbox-update@{{INSTANCE}}.service

[Timer]
# Run daily at 02:00
OnCalendar=*-*-* 02:00:00
# Catch up on missed runs after boot
Persistent=true
AccuracySec=1min
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
